# infra/crontab/daily.sh
#!/bin/bash

(cd /var/gg/gloriousgeckos/backend && npm run cronjob:daily)